# dailycron.cron
#Min	Hour	Day	Month	Wday	script
*	*	*	*	*	/home/uwslowcontrol/pi_db/checkpi_db.sh
#*	*	*	*	*	pgrep -f pi_db.py || python ~/pi_db/pi_db.py

