# cronjobs/ios/dailycron.cron
#Min	Hour	Day	Month	Wday	script
26	11	*	*	*	/home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/diskUsage.sh 
26	12	*	*	Tue	/home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/compact.sh
26	23	*	*	Mon	python /home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/deleteOldDocs.py
*	*	*	*	*	python /home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/checkCouchReplication.py		
*	*	*	*	*	/home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/checkPoll.sh
*	*	*	*	*	/home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/check1min.sh
*/15	*	*	*	*	/home/slowcontroller/SNOPlusSlowControl/cronjobs/ios/check15min.sh 
26      1       *       *       *       /etc/init.d/ntpd start
