# cronjobs/pi_db/dailycron.cron
#Min	Hour	Day	Month	Wday	script
*	*	*	*	*	/home/uwslowcontrol/SNOPlusSlowControl/SNOPlusSlowControl/cronjobs/pi_db/checkpi_db.sh
*	*	*	*	*	/home/uwslowcontrol/SNOPlusSlowControl/SNOPlusSlowControl/cronjobs/pi_db/checkctemp.sh

